# myconfig/systemd/app-misc/unison/user/unison@.timer
[Unit]
Description=Unison
StopWhenUnneeded=yes

[Timer]
OnActiveSec=30min
OnUnitInactiveSec=30min

[Install]
WantedBy=timers.target