# myconfig/systemd/sys-apps/systemd/user/environment.target
[Unit]
Description=Environment
Requires=dbus.socket
Before=basic.target
IgnoreOnIsolate=true

[Install]
WantedBy=default.target
