# myconfig/systemd/x11-server/xorg-server/user/X@.target
[Unit]
Description=sync-point for X (Xorg, Xwayland) target at DISPLAY :%i
Requires=basic.target
RefuseManualStart=yes
StopWhenUnneeded=yes
