# myconfig/systemd/x11-server/xorg-server/user/xorg@.socket
[Unit]
Description=Socket for X.org at display %i

[Socket]
ListenStream=/tmp/.X11-unix/X%i

