# myconfig/systemd/x11-server/xorg-server/user/xorg.service
[Unit]
Description=The X server for X.org

[Service]
ExecStart=/usr/bin/Xorg ${DISPLAY} vt7

[Install]
WantedBy=xserver.target
