# myconfig/systemd/x11-server/xorg-server/user/xorg.socket
[Unit]
Description=Socket for X.org

[Socket]
ListenStream=/tmp/.X11-unix/X0
ExecStartPost=-/usr/bin/systemctl --user set-environment DISPLAY=:0


