# myconfig/systemd/inputmethods/fcitx/user/fcitx.service
[Unit]
Description=Fcitx Input Method

[Service]
ExecStart=/usr/bin/fcitx -D -r

