# myconfig/systemd/x11-apps/rxvt-unicode/user/urxvtd.service
[Unit]
Description=rxvt-unicode daemon

[Service]
ExecStart=/usr/bin/urxvtd

[Install]
WantedBy=graphical.target
