# myconfig/systemd/x11-apps/rxvt-unicode/user/urxvtd.service
[Unit]
Description=Urxvt Terminal Daemon
Requires=urxvtd.socket

[Service]
ExecStart=/usr/bin/urxvtd -q
