# myconfig/systemd/app-crypt/gnupg/user/gpg-agent.service
[Unit]
Description=GnuPG private key agent
IgnoreOnIsolate=true

[Service]
Environment=GNUPGHOME=%h/.config/gnupg
Type=forking
ExecStart=/usr/bin/gpg-agent --daemon
ExecStartPost=-/usr/bin/systemctl --user import-environment GNUPGHOME
Restart=on-abort

[Install]
WantedBy=default.target
